# .config/systemd/user/animated-wallpaper.service
[Unit]
Description=Animated Wallpaper

[Service]
ExecStart=video-bg %h/.config/wallpapers/default-animated

[Install]
WantedBy=default.target
