# .config/systemd/user/compton@.service
[Unit]
Description=Compton Compositor

[Service]
Environment=DISPLAY=:%i
ExecStart=compton

[Install]
WantedBy=default.target
