# .config/systemd/user/compton.service
[Unit]
Description=Compton Compositor

[Service]
ExecStart=compton

[Install]
WantedBy=default.target
