# .config/systemd/user/dunst@.service
[Unit]
Description=Dunst Notification Daemon

[Service]
Environment=DISPLAY=:%i
ExecStart=dunst

[Install]
WantedBy=default.target
