# .config/systemd/user/dunst.service
[Unit]
Description=Dunst Notification Daemon

[Service]
ExecStart=dunst

[Install]
WantedBy=default.target
