# .config/systemd/user/network-manager-applet@.service
[Unit]
Description=Applet for managing network connections

[Service]
Environment=DISPLAY=:%i
ExecStart=nm-applet

[Install]
WantedBy=default.target
