# .config/systemd/user/polybar-wss.service
[Unit]
Description=Polybar wss bar

[Service]
ExecStart=polybar wss

[Install]
WantedBy=default.target
