# .config/systemd/user/polybar-wss@.service
[Unit]
Description=Polybar wss bar

[Service]
Environment=DISPLAY=:%i
ExecStart=polybar wss

[Install]
WantedBy=default.target
