# .config/systemd/user/polybar-main.service
[Unit]
Description=Polybar main bar

[Service]
ExecStart=polybar main

[Install]
WantedBy=default.target
