# .config/systemd/user/polybar-main@.service
[Unit]
Description=Polybar main bar

[Service]
Environment=DISPLAY=:%i
ExecStart=polybar main

[Install]
WantedBy=default.target
