# .config/systemd/user/slack@.service
[Unit]
Description=Slack

[Service]
Environment=DISPLAY=:%i
ExecStart=slack

[Install]
WantedBy=default.target
