# .config/systemd/user/udiskie.service
[Unit]
Description=Removable disk automounter with tray icon

[Service]
ExecStart=udiskie -s

[Install]
WantedBy=default.target
