# .config/systemd/user/waybar-main.service
[Unit]
Description=Waybar main bar

[Service]
ExecStart=waybar

[Install]
WantedBy=default.target
