# .config/systemd/user/anbox-session-manager.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/anbox session-manager --gles-driver=host
